# com.stargazer.dailythreadbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.stargazer.dailythreadbot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <!-- Replace this with the absolute path to run-daily-thread.sh before loading the plist. -->
        <string>/absolute/path/to/daily-thread-bot/run-daily-thread.sh</string>
    </array>

    <!-- Example: Mon-Fri at 9:30 AM local machine time. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>6</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>10</integer><key>Minute</key><integer>6</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/absolute/path/to/daily-thread-bot/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/absolute/path/to/daily-thread-bot/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
